# service-files/ai-backend.service
[Unit]
Description=AI Backend FastAPI Service
After=network.target

[Service]
User=djr
Group=djr
WorkingDirectory=/home/djr/ai-server/ai-backend
ExecStart=/bin/bash /home/djr/ai-server/ai-backend/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

