# service-files/ollama-warmup.service
[Unit]
Description=Warm up Ollama model on boot
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/djr/ai-server/ollama-warmup/ollama-warmup.sh

[Install]
WantedBy=multi-user.target

